# psmrc.service
[Unit]
Description=PSMRC-Nodles - Minecraft to PS Vita resource pack converter
After=network.target

[Service]
Type=simple
User=psmrc
Group=psmrc
WorkingDirectory=/opt/psmrc-nodles
ExecStart=/usr/bin/node index.js # If this is not your Node.js path replace it with a correct path that 
Restart=on-failure
RestartSec=5

Environment=NODE_ENV=production
Environment=PORT=3001 # If occupied use another free port
Environment=HOST=127.0.0.1

# Basic hardening - safe defaults for an app that only needs its own directory
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
